# deploy/escra.service
[Unit]
Description=Escra — OLX Lviv scraper + Telegram Web App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Running as root (matches how the repo was cloned on the droplet: /root/escra).
User=root
WorkingDirectory=/root/escra
EnvironmentFile=/root/escra/.env
# uv creates .venv on `uv sync`. Point straight at the venv interpreter.
ExecStart=/root/escra/.venv/bin/python -m app.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
